# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf

[progrem:cron]
command=/bin/bash -c "cron &"
autostart=true
autorestart=false
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

# [program:execute_shell_script]
# user=iptvboss
# command=/headless/entrypoint.sh
# autostart=true
# autorestart=false
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0

[program:xcserver]
command=/bin/bash -c "/usr/bin/iptvboss -xcserver"
startretries=2
autorestart=false
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

# [program:vncserver]
# command=/bin/bash -c "exec /dockerstartup/vnc_startup.sh --wait"
# startretries=5
# autorestart=false
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
